# ops/backup/lembaga-weekly-whatsapp-restart.timer
[Unit]
Description=Jadwal restart WhatsApp bot setiap Minggu pukul 01.00 WIB

[Timer]
OnCalendar=Sun *-*-* 01:00:00 Asia/Jakarta
Persistent=false
AccuracySec=1min
RandomizedDelaySec=0
Unit=lembaga-weekly-whatsapp-restart.service

[Install]
WantedBy=timers.target
